# linux_net_agent/systemd/linux-net-agent.service
[Unit]
Description=Linux-Net C++ host agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linux-net
Group=linux-net
ExecStart=/usr/local/bin/linux-net-agent --host 0.0.0.0 --port 7070
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
